# src/assets/etc/systemd/system/reflector.service.d/override.conf
[Service]
ExecStart=
ExecStart=/usr/bin/reflector --save /etc/pacman.d/mirrorlist --protocol https --sort rate
