# src/assets/etc/systemd/system/getty@tty7.service.d/override.conf
[Service]
ExecStart=
ExecStart=-/usr/bin/agetty --skip-login --nonewline --noissue --autologin caretakr --noclear %I $TERM
ExecStartPost=/bin/chvt 7
